# bash_history_history/backup_bash_history.service
[Unit]
Description=Backup ~/.bash_history

[Service]
Type=simple
ExecStart=bash %h/tools/commit_bash_history.sh