# bash_history_history/backup_bash_history.timer
[Unit]
Description=Trigger backing up ~/.bash_history

[Timer]
OnCalendar=weekly
Unit=backup_bash_history.service

[Install]
WantedBy=timers.target