# backup_shell_history/backup_bash_history.service
[Unit]
Description=Backup ~/.bash_history

[Service]
Type=simple
ExecStart=%h/.local/bin/commit_bash_history.sh %h/.bash_history_history
